# config/crontab
0 */8 * * * screen -S EnterpriseApp -dm /bin/bash /home/master/enterprise/bin/run_production
